# docker/crontab
# Crontab for database backups
# This file is used to schedule automatic backups

# Environment variables
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Backup schedule - Daily at 2 AM
0 2 * * * /backup.sh >> /var/log/backup.log 2>&1

# Clean old logs - Weekly on Sunday at 3 AM
0 3 * * 0 find /var/log -name "*.log" -type f -mtime +7 -delete

# Health check - Every 5 minutes
*/5 * * * * curl -f http://app:3000/api/health || echo "Health check failed at $(date)" >> /var/log/health.log